# unit_backups/overlay-updater.timer
[Unit]
Description=Run overlay-updater every 30 seconds

[Timer]
OnBootSec=15s
OnUnitActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target
